# rhel6-jenkins-shared-slave/supervisord.conf
[supervisord]  
nodaemon=true  
  
[program:sshd]  
command=/sbin/service sshd start 
